# files/usr/lib/systemd/user/azure-distrobox.service
[Unit]
Description=Update distrobox for current user
Wants=network-online.target
After=network-online.service

[Service]
Type=simple
ExecStart=/usr/libexec/azure-distrobox
Restart=on-failure
RestartSec=30
StartLimitInterval=0

[Install]
WantedBy=default.target